# templates/chatops-restart.cron.j2
MAILTO=""
*/10 * * * * root systemctl restart chatops.service >> {{ lookup('ini', 'LOG_FILE section=chatops-restart.cron file=chatops_env.config') }}